# ngrok.service
[Unit]
Description=ngrok tunnel service
After=network.target

[Service]
Type=simple
User=fd1337
Group=fd1337
WorkingDirectory=/home/fd1337
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/ngrok http 3000 --log=stdout
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
